# webui/deploy/okforge-webui-standalone.service
# Standalone okforge-webui: uvicorn serves the frontend, the API, and the
# MCP server itself — no Apache. (deploy.sh installs the Apache-fronted
# variant instead; use one or the other.)
#
# Install by hand:
#   sudo cp webui/deploy/okforge-webui-standalone.service \
#       /etc/systemd/system/okforge-webui.service
#   ... edit User=/Group= and paths below ...
#   sudo systemctl daemon-reload && sudo systemctl enable --now okforge-webui

[Unit]
Description=okforge web UI (standalone uvicorn, python -m webui)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=okforge
Group=okforge
WorkingDirectory=/opt/okforge/okforge-webui
ExecStart=/opt/okforge/okforge-webui/.venv/bin/python -m webui
Restart=on-failure
RestartSec=5
# Bind + backend config (see webui/config.py for the full list), e.g.:
# Environment=OPENKB_WEBUI_HOST=0.0.0.0
# Environment=OPENKB_WEBUI_PORT=8500
# Environment=OPENKB_WEBUI_ENDPOINTS=gpu1=http://gpu1:8080/v1
# Environment=OPENKB_WEBUI_MD_OUT=/opt/okforge/md-out
# Environment=OPENKB_WEBUI_TRASH=/opt/okforge/trash

[Install]
WantedBy=multi-user.target
